# myservice.service
[Unit]
Description=Service to download the html changes

[Service]
Type=oneshot
ExecStart=/opt/htmlDownloader/service.sh
User=paulograbin